# xrootd/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes=1024MB
logfile_backups=2
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file = /tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:xrd1]
command=xrootd -R xrootd -n server1 -c /etc/xrootd/xrootd-server1.cfg
name=xrd1

[program:xrd2]
command=xrootd -R xrootd -n server2 -c /etc/xrootd/xrootd-server2.cfg
name=xrd2

[program:xrd3]
command=xrootd -R xrootd -n server3 -c /etc/xrootd/xrootd-server3.cfg
name=xrd2
